# pi/systemd/picam-capture.timer
[Unit]
Description=PiCam capture gate — trigger every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
