# pi/systemd/picam-capture.service
[Unit]
Description=PiCam capture gate
After=network-online.target

[Service]
Type=oneshot
RuntimeDirectory=picam
RuntimeDirectoryPreserve=yes
ExecStart=/usr/local/bin/capture-gate.sh
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
